# heartbeat/com.celticxfer.moltbook.heartbeat.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.celticxfer.moltbook.heartbeat</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Volumes/FS001/pythonscripts/moltbot/heartbeat/celticxfer_heartbeat.sh</string>
    </array>
    <!-- 14400 seconds = 4 hours -->
    <key>StartInterval</key>
    <integer>14400</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Volumes/FS001/pythonscripts/moltbot/heartbeat/heartbeat.log</string>
    <key>StandardErrorPath</key>
    <string>/Volumes/FS001/pythonscripts/moltbot/heartbeat/heartbeat.error.log</string>
    <!-- Don't run if the external volume isn't mounted -->
    <key>WatchPaths</key>
    <array>
        <string>/Volumes/FS001</string>
    </array>
</dict>
</plist>
